# ansible/files/mixmaster-inbox-watcher.path
[Unit]
Description=Mixmaster inbox path watcher

[Path]
DirectoryNotEmpty=/home/mixmaster/Builds/INBOX
MakeDirectory=true
DirectoryMode=0700

[Install]
WantedBy=multi-user.target

# Local Variables:
# mode: conf
# End:
