# services/mixmaster-bridge.socket
[Unit]
Description=Mixmaster Bridge Socket

[Socket]
ListenStream=8585
Accept=yes
ReusePort=true

[Install]
WantedBy = sockets.target
