# systemd/mixmaster-ondemand-bridge.socket
[Unit]
Description=Mixmaster Socket Bridge for on-demand

[Socket]
ListenStream=8586
Accept=yes
ReusePort=true

[Install]
WantedBy = sockets.target
