# systemd/mixmaster-inbox-watcher.path
[Unit]
Description=Mixmaster inbox path watcher

[Path]
DirectoryNotEmpty=/home/%u/Documents/mixmaster/INBOX
MakeDirectory=true
DirectoryMode=0700

[Install]
WantedBy=multi-user.target
